# tools/com.decosta.git-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.decosta.git-push</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zhutovoleg/Doc/CLAUDE/CLAUDE_CODE/tools/auto_git_push.sh</string>
    </array>

    <!-- Запуск каждый день в 12:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/zhutovoleg/Doc/CLAUDE/CLAUDE_CODE/tools/git_push.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zhutovoleg/Doc/CLAUDE/CLAUDE_CODE/tools/git_push.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
